# agent-python/dist/macos/payload/Library/LaunchAgents/com.paccontrol.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.paccontrol.agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/pac-control-agent/PACControlAgent</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PAC_AGENT_CONFIG</key>
    <string>/opt/pac-control-agent/config.json</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/pac-control-agent.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/pac-control-agent.err.log</string>
</dict>
</plist>
